# deploy/systemd/webrtc-camera-share.service.example
[Unit]
Description=WebRTC Camera Share
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=webrtc-camera-share
Group=webrtc-camera-share
WorkingDirectory=/var/lib/webrtc-camera-share
EnvironmentFile=/etc/webrtc-camera-share/app.env
ExecStart=/usr/local/bin/webrtc-camera-share-server
Restart=on-failure
RestartSec=3s
TimeoutStartSec=20s
TimeoutStopSec=20s
KillSignal=SIGTERM

NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
LockPersonality=true
MemoryDenyWriteExecute=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictAddressFamilies=AF_INET AF_INET6

LimitNOFILE=65536
TasksMax=256
MemoryMax=512M

[Install]
WantedBy=multi-user.target
